# systemd/dot-config/systemd/user/update-ssh-config-incus.service
[Unit]
Description=Regenerate SSH config from Incus instances

[Service]
ExecStart=/usr/local/bin/generate_incus_ssh
