# systemd/dot-config/systemd/user/update-ssh-config-incus.path
[Unit]
Description=Trigger SSH config regeneration on Incus changes

[Path]
PathChanged=/var/lib/incus
PathChanged=/run/incus
Unit=update-ssh-incus.service

[Install]
WantedBy=multi-user.target
